# HID/payload-ui.service
[Unit]
Description=QML Payload UI
Conflicts=getty@tty2.service
After=systemd-user-sessions.service

[Service]
Environment="QT_QPA_EVDEV_KEYBOARD_PARAMETERS=grab=0"
Environment="QT_QPA_EVDEV_TOUCH_PARAMETERS=/dev/input/event2"
ExecStart=/usr/local/bin/payload-app -platform linuxfb:fb=/dev/fb0 -plugin evdevtouch:/dev/input/event3
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty2
Restart=always
User=root

[Install]
WantedBy=multi-user.target