# HID/payload-ui.service
[Unit]
Description=HID UI
Conflicts=getty@tty2.service
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/cage -s -- /usr/local/bin/payload-app
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty2
Restart=always
User=root
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="WAYLAND_DISPLAY=wayland-1"

[Install]
WantedBy=multi-user.target